# defaults/media-centarr-backend.service
[Unit]
Description=Media Centarr Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
# Set this to your media-centarr release directory, e.g.:
#   ExecStart=%h/.local/lib/media-centarr/bin/media_centarr start
ExecStart=/path/to/media_centarr start
ExecStop=/path/to/media_centarr stop
Restart=on-failure
RestartSec=5
KillMode=process
TimeoutStopSec=30
LimitNOFILE=65536

[Install]
WantedBy=default.target
