# defaults/media-centaur-backend-dev.service
[Unit]
Description=Media Centaur Backend (Dev)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
# Set this to your media-centaur backend source directory, e.g.:
#   WorkingDirectory=%h/src/media-centaur/backend
WorkingDirectory=/path/to/media-centaur/backend
ExecStart=/usr/bin/mise exec -- mix phx.server
Environment=MIX_ENV=dev
Environment=MIX_OS_DEPS_COMPILE_PARTITION_COUNT=8
Environment=ELIXIR_ERL_OPTIONS="-name media_centaur_dev@127.0.0.1"
Restart=on-failure
RestartSec=5
KillMode=process
TimeoutStopSec=30
LimitNOFILE=65536

[Install]
WantedBy=default.target
